# docker-test/build-plot/crontab
* * * * * root /var/local/getMetadata.sh
* * * * * root /var/local/makePlot.sh
# new line
